# etc/sv/acpid/run
#!/bin/sh
exec acpid -d -a /etc/acpi/acpid.conf -M /etc/acpi/acpi.map >/dev/null
