# etc/sv/dhcpcd/run
#!/bin/sh
exec dhcpcd -BMqq 2>&1 >/dev/null
